# packages/nvidia-fabrics-manager/nvidia-fabrics-manager.service
[Unit]
Description=Start NVIDIA fabrics manager
RefuseManualStart=true
RefuseManualStop=true
After=kubelet.service
Wants=kubelet.service

[Service]
ExecStart=/usr/bin/nvidia-fabric-manager-devel -p /sbin/ldconfig
Type=simple
TimeoutSec=0
RestartSec=2
Restart=always
StandardError=journal+console

[Install]
WantedBy=multi-user.target
